# lxc/gitea/gitea-check-config.timer
[Unit]
Description=Run Gitea configuration check daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
