# lxc/tapmap/etc_geoipupdate/geoipupdate.service.d/override.conf
[Service]
# Allow access to /root/.local/share/TapMap for database storage
ProtectHome=no
# Add TapMap data directory as writable
ReadWritePaths=/root/.local/share/TapMap/
